# ansible/roles/rustsrv/templates/lupdator.service.j2
[Unit]
Description=Run updator every 20 minutes

[Service]
Type=simple
WorkingDirectory={{ ldest_dir }}
ExecStart={{ ldest_dir }}/loong-updator
