# ansible/roles/dbserver/files/backuprepo.service
[Unit]
Description=Borg backup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backuprepo
